# debian/buddydrive.service
[Unit]
Description=BuddyDrive P2P Folder Sync
Documentation=https://github.com/anomalyco/buddydrive
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/buddydrive start
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Environment - separate config and data directories
Environment="BUDDYDRIVE_CONFIG_DIR=/etc/buddydrive"
Environment="BUDDYDRIVE_DATA_DIR=/var/lib/buddydrive"

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/buddydrive

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
